# root/etc/s6-overlay/s6-rc.d/base-apply-permissions/up
/etc/s6-overlay/scripts/base-apply-permissions
