# releases/noble/root/etc/s6-overlay/s6-rc.d/base-set-up-user/up
/etc/s6-overlay/s6-rc.d/base-set-up-user/run
