# test/root/etc/s6-overlay/s6-rc.d/test/run
#!/command/with-contenv bash

set -e
exec /etc/s6-overlay/scripts/test
